# deploy/catch-missed-alerts-detect.service
[Unit]
Description=Catch missed critical Slack alerts — hourly detect pass
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/DetectUnackedAlerts
ExecStart=/home/ec2-user/DetectUnackedAlerts/venv/bin/python /home/ec2-user/DetectUnackedAlerts/alert_watcher.py detect
# don't let one bad run wedge the timer
TimeoutStartSec=120
