# deploy/catch-missed-alerts-check.service
[Unit]
Description=Catch missed critical Slack alerts — 5-min ack re-check / re-ping
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/home/ec2-user/DetectUnackedAlerts
ExecStart=/home/ec2-user/DetectUnackedAlerts/venv/bin/python /home/ec2-user/DetectUnackedAlerts/alert_watcher.py check
TimeoutStartSec=120
