# deploy/catch-missed-alerts-detect.service
[Unit]
Description=Catch missed critical Slack alerts — hourly detect pass
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=ec2-user
WorkingDirectory=/opt/catch_missed_alerts
ExecStart=/opt/catch_missed_alerts/venv/bin/python /opt/catch_missed_alerts/alert_watcher.py detect
# don't let one bad run wedge the timer
TimeoutStartSec=120
