# deploy/catch-missed-alerts-escalate.timer
[Unit]
Description=Run catch-missed-alerts escalate once a day at 07:00 (end of night window)

[Timer]
# 07:00 local time daily — validates the overnight (21:00–07:00) criticals
OnCalendar=*-*-* 07:00:00
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
